# ops/sports-predictions.service
[Unit]
Description=Sports Predictions application
After=network-online.target
[Service]
Type=simple
ExecStart=/home/sports-predictions/start.sh
Restart=on-failure
[Install]
WantedBy=multi-user.target
